# .config/systemd/user/system-backup.timer
[Unit]
Description=Backup system settings

[Timer]
OnCalendar=Mon,Tue,Wed,Thu,Fri 08:00

[Install]
WantedBy=timers.target
