# .config/systemd/user/wakeonlan@.service
[Unit]
Description=Send magick packets

[Service]
Type=oneshot
ExecStart=/usr/bin/wakeonlan %i

[Install]
WantedBy=default.target
